# systemd/tele-quant-weekly.timer
[Unit]
Description=Run Tele Quant weekly digest every Saturday at 09:00

[Timer]
OnCalendar=Sun *-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
