# systemd/tele-quant-backlog-us.timer
[Unit]
Description=Run Tele Quant Backlog Refresh US every weekday at 07:00 UTC

[Timer]
OnCalendar=Mon..Fri *-*-* 07:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
