# systemd/tele-quant-weekend-macro.service
[Unit]
Description=Tele Quant weekend macro-only digest (no stock analysis)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/quant_spillover/tele_quant
Environment=DIGEST_MODE=fast
ExecStart=%h/.local/bin/uv run tele-quant once --send --hours 4 --macro-only
