# systemd/tele-quant-daily-alpha-us.timer
[Unit]
Description=Run Tele Quant Daily Alpha Picks US every weekday at 22:00 KST (13:00 UTC)

[Timer]
OnCalendar=Mon..Fri *-*-* 13:00:00
Persistent=true

[Install]
WantedBy=timers.target
