# systemd/tele-quant-backlog-kr.service
[Unit]
Description=Tele Quant Backlog Refresh — KR 07:30 KST (weekdays)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/quant_spillover/tele_quant
ExecStart=%h/.local/bin/uv run tele-quant backlog-refresh --market KR --days 30 --save
