# systemd/tele-quant-daily-alpha-kr-intraday-10.timer
[Unit]
Description=Daily Alpha KR 장중 10:00 KST 스캔

[Timer]
OnCalendar=Mon..Fri *-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
