# systemd/tele-quant-backlog-report.timer
[Unit]
Description=Run Tele Quant Backlog Report every Sunday at 22:00 KST

[Timer]
OnCalendar=Sun *-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
